# airootfs/etc/systemd/system/load-config-before-xorg.service
[Unit]
Description=Load configs before xorg starts

[Service]   
Type=oneshot
ExecStart=/root/pre-xorg.sh

[Install]
WantedBy=multi-user.target




